# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/guild_relationships_2016_npr/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/guild_relationships_2016_npr/crontab.log 2>&1
